# deploy_staging/etc/systemd/system/dispatch-watchdog.timer.d/oncalendar.conf
# FALA1 watchdog-close (2026-07-02) — durability fix for audyt 2.0 finding A.
#
# Problem: base timer scheduled ONLY via OnUnitActiveSec=4h, with Persistent=true that
# is a NO-OP (Persistent= only affects OnCalendar= timers — systemd.timer(5)). On
# 26.06 19:43 a collective `disable` left the cron-health watchdog silent for 5 days;
# the stale-detector was itself in the dead set, so nothing noticed.
#
# Fix: anchor to a wall-clock OnCalendar (every 4h) so the schedule survives
# daemon-reload / re-enable, and Persistent=true now actually catches up a missed run
# after downtime.
#
# No double-fire: assigning "" to a timer option resets the ENTIRE timer list
# (monotonic + OnCalendar — systemd.timer(5)); we then re-declare exactly one basis.
# Expected merged result — verify with: systemctl cat dispatch-watchdog.timer
#   OnBootSec=15min ; OnCalendar=*-*-* 00/04:00:00 ; Persistent=true ; RandomizedDelaySec=5min
# systemd-analyze calendar '*-*-* 00/04:00:00'  -> 00:00,04:00,08:00,12:00,16:00,20:00
[Timer]
OnUnitActiveSec=
OnBootSec=15min
OnCalendar=*-*-* 00/04:00:00
Persistent=true
