# deploy_staging/etc/systemd/system/dispatch-downstream-crosscheck.service.d/cron_health_success.conf
# FALA1 watchdog-close (2026-07-02) — audyt 2.0: false-"failed" oneshot.
# Runs every 5 min and exits 0 (verified: Jul 02 06:38 "downstream_status=ok /
# Deactivated successfully") but had no success recorder → ledger frozen at
# failed/last_success=null since a one-off timeout 2026-06-13. ExecStopPost records
# success so the ledger clears and the 1.0h stall net (registry backfill) has an anchor.
[Service]
ExecStopPost=/root/.openclaw/workspace/scripts/dispatch_v2/observability/record_oneshot_success.sh %n
