# deploy_staging/etc/systemd/system/dispatch-delivered-integrity.timer.d/oncalendar.conf
# FALA1 watchdog-close (2026-07-02) — durability fix for audyt 2.0 finding A / §4.
#
# dispatch-delivered-integrity.timer was in the collective disable of 26.06 19:43:41
# (delivered_at=None integrity monitor, B3/B5). Base scheduled only via
# OnUnitActiveSec=20min → no wall-clock anchor. Re-anchor to OnCalendar so a
# daemon-reload / collective disable cannot silently stop it again.
#
# No double-fire: "" resets the whole timer list (systemd.timer(5)); re-declare one basis.
# Persistent left at base default (false) — 20-min cadence needs no catch-up burst.
# Expected merged — verify: systemctl cat dispatch-delivered-integrity.timer
#   OnBootSec=5min ; OnCalendar=*:0/20 ; AccuracySec=30s ; Persistent=false
# systemd-analyze calendar '*:0/20'  -> every 20 min on the wall clock (:00,:20,:40)
[Timer]
OnUnitActiveSec=
OnBootSec=5min
OnCalendar=*:0/20
