# deploy_staging/etc/systemd/system/dispatch-retro-learning.service.d/cron_health_success.conf
# FALA1 watchdog-close (2026-07-02) — audyt 2.0: false-"failed" oneshot.
# Runs daily 04:30 and exits 0 (verified: Jul 02 04:30 "Deactivated successfully") but
# had no success recorder → ledger frozen at failed/last_success=null since 2026-06-08.
# ExecStopPost records success so the ledger clears and the 25h stale net (registry
# backfill) has a real anchor.
[Service]
ExecStopPost=/root/.openclaw/workspace/scripts/dispatch_v2/observability/record_oneshot_success.sh %n
