# deploy_staging/etc/systemd/system/dispatch-retro-learning.service.d/onfailure.conf
# FALA1 watchdog-close (2026-07-02) — parity: retro-learning had NO OnFailure drop-in.
# Without it a real failure records nothing to the ledger (the ExecStopPost recorder
# only fires on success), so a failing retro run would be doubly silent. 1:1 with the
# standard dispatch-* OnFailure handler.
[Unit]
OnFailure=dispatch-onfailure-alert@%n.service
