# deploy_staging/etc/dispatch-bundle-calib-shadow.service.d/lambda0-recollect.conf
# Re-collect λ=0 (2026-07-03, checklist bug4-logger_raport §4 / o2-capz_raport §5).
# Cel: korpus ground-truth O2 BEZ λ-zawyżenia (λ siedzi w SELEKCJI CALIB — post-hoc
# przeliczenie niemożliwe, tylko re-collect). Zbierać ≥2-3 dni → bundle_calib_review
# z BUNDLE_CALIB_CORPUS=...l0.jsonl → porównać engine-improved% z 7.3% (λ=1.5).
# Rozłączne pliki output+state = zero skażenia korpusu λ=1.5 (zostaje nietknięty).
# Rollback: usuń ten plik + daemon-reload (wraca λ=1.5 i stare ścieżki).
[Service]
Environment=BUNDLE_CALIB_LAMBDA_CZAS=0
Environment=BUNDLE_CALIB_OUT_JSONL=/root/.openclaw/workspace/dispatch_state/bundle_calib_shadow_l0.jsonl
Environment=BUNDLE_CALIB_STATE_PATH=/root/.openclaw/workspace/dispatch_state/bundle_calib_shadow_state_l0.json
