# deploy_staging/etc/systemd/system/dispatch-cod-weekly.service.d/cron_health_success.conf
# FALA1 watchdog-close (2026-07-02) — audyt 2.0 finding B (parity with siblings).
# Refresh cron_health last_success ONLY on SERVICE_RESULT=success (helper checks it),
# so once the Monday run succeeds (week block present) the ledger clears and the 192h
# staleness net has a real anchor. 1:1 with dispatch-cod-weekly-preflight.service.d.
[Service]
ExecStopPost=/root/.openclaw/workspace/scripts/dispatch_v2/observability/record_oneshot_success.sh %n
