# deploy_staging/etc/systemd/system/dispatch-cod-weekly.service.d/onfailure.conf
# FALA1 watchdog-close (2026-07-02) — audyt 2.0 finding B.
# dispatch-cod-weekly (COD payouts, Mon 08:00 Warsaw) was the ONLY dispatch job with
# NO .service.d/ at all → no OnFailure → its weekly exit-1 (Jun 08/15/22/29, missing
# week block in the sheet) was silent at the systemd layer. 1:1 with the drop-in every
# other dispatch-* has (e.g. dispatch-cod-weekly-preflight.service.d/onfailure.conf).
[Unit]
OnFailure=dispatch-onfailure-alert@%n.service
