# deploy_staging/etc/systemd/system/dispatch-cod-weekly.service.d/resource_limits.conf
# FALA1 watchdog-close (2026-07-02) — bliźniak parity (audyt 2.0 §3 fix proposal #1).
# cod-weekly does the same gspread batch work as its twins; give it the identical
# oneshot-batch resource profile (1:1 with dispatch-cod-weekly-preflight.service.d /
# dispatch-daily-accounting.service.d). Prevents an OOM/hang from being unbounded.
[Service]
MemoryHigh=400M
MemoryMax=600M
CPUQuota=100%
TasksMax=256
TimeoutStartSec=600
OOMScoreAdjust=100
