# systemd/dispatch-reassignment-shadow-review.timer
[Unit]
Description=Ziomek reassignment-shadow REVIEW — jednorazowo 2026-06-23 19:00 UTC (= 21:00 Warsaw, po obu peakach)

[Timer]
OnCalendar=2026-06-23 19:00:00 UTC
Persistent=true
AccuracySec=1min
Unit=dispatch-reassignment-shadow-review.service

[Install]
WantedBy=timers.target
