# systemd/dispatch-pickup-lateness-review.timer
[Unit]
Description=Ziomek pickup-lateness REVIEW — jednorazowo 2026-06-23 18:00 UTC (= 20:00 Warsaw, po peaku)

[Timer]
OnCalendar=2026-06-23 18:00:00 UTC
Persistent=true
AccuracySec=1min
Unit=dispatch-pickup-lateness-review.service

[Install]
WantedBy=timers.target
