# systemd/dispatch-b3-trial-monitor.service
# B3-MONITOR — dzienny auto-monitor trialu ENABLE_NO_GPS_UNCERTAINTY_PENALTY
# z auto-rollbackiem przy breachu (override>50% lub on-time<60% przy n>=15).
# Auto-rollback celuje TYLKO w stan bezpieczny (flaga=false = powrót do demote V3.16).
# Tryb alert-only: B3_MONITOR_AUTO_ROLLBACK=0 (nie tyka flagi, tylko alert).
#
# INSTALL (po przeglądzie team-lead):
#   cp dispatch_v2/systemd/dispatch-b3-trial-monitor.{service,timer} /etc/systemd/system/
#   mkdir -p /etc/systemd/system/dispatch-b3-trial-monitor.service.d
#   cp dispatch_v2/systemd/onfailure.conf /etc/systemd/system/dispatch-b3-trial-monitor.service.d/
#   systemctl daemon-reload && systemctl enable --now dispatch-b3-trial-monitor.timer
[Unit]
Description=Ziomek B3 trial monitor — dzienny audyt no_gps uncertainty + auto-rollback
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/.openclaw/workspace/scripts
# Domyślnie auto-rollback ON. Aby tylko alertować bez zmiany flagi: Environment=B3_MONITOR_AUTO_ROLLBACK=0
Environment=B3_MONITOR_AUTO_ROLLBACK=1
Environment=B3_FLIP_DATE=2026-06-20
ExecStart=/root/.openclaw/venvs/dispatch/bin/python -m dispatch_v2.tools.b3_trial_monitor
StandardOutput=append:/root/.openclaw/workspace/scripts/logs/b3_trial_monitor.log
StandardError=append:/root/.openclaw/workspace/scripts/logs/b3_trial_monitor.log
# Sandbox / hardening (wzór: assistant-rag-index / dispatch resource_limits)
MemoryMax=400M
MemoryHigh=320M
CPUQuota=50%
TasksMax=128
NoNewPrivileges=true
PrivateTmp=true
OOMScoreAdjust=-100
