# systemd/dispatch-plan-recheck.timer
[Unit]
Description=Ziomek V3.19c Plan Recheck — every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=false
Unit=dispatch-plan-recheck.service

[Install]
WantedBy=timers.target
