# systemd/dispatch-pending-resweep-shadow.service
[Unit]
Description=Ziomek pending-resweep SHADOW — globalny re-ranking WISZĄCYCH propozycji co 1 min; zapisuje PII-free JSONL i atomowy stan kontrfaktyczny G6. Zero mutacji panel/Telegram/pending/orders. ENABLE_PENDING_RESWEEP default OFF=no-op; PENDING_RESWEEP_LIVE default OFF.
Documentation=file:/root/.openclaw/workspace/scripts/dispatch_v2/tools/pending_global_resweep.py
Documentation=file:/root/.openclaw/workspace/scripts/dispatch_v2/eod_drafts/2026-06-24/pending_global_resweep.md
After=network.target
OnFailure=dispatch-onfailure-alert@%n.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/.openclaw/workspace/scripts
ExecStart=/root/.openclaw/venvs/dispatch/bin/python -m dispatch_v2.tools.pending_global_resweep
# Peak-friendly: niski priorytet CPU + I/O w tle (assess_order woła OSRM; nie wypierać hot-pathu)
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Bezpiecznik MAX_HANGING=8 → O(K^2) assess; oneshot, nie powinno przekroczyć minuty
TimeoutStartSec=90
# Sandbox / pamięć — proces mały; zapis ograniczony logicznie do JSONL i stanu G6.
MemoryHigh=300M
MemoryMax=400M
TasksMax=256
# Shadow = preferowany do ubicia gdy host OOM (ginie PIERWSZY, NIE wypiera dispatch-shadow/panel-watcher)
OOMScoreAdjust=200
StandardOutput=append:/root/.openclaw/workspace/scripts/logs/pending_global_resweep.log
StandardError=append:/root/.openclaw/workspace/scripts/logs/pending_global_resweep.log

[Install]
WantedBy=multi-user.target
