# systemd/dispatch-reassignment-shadow-eval.service
[Unit]
Description=Ziomek reassignment-shadow EVAL — jednorazowy materiał GO/NO-GO (would_reassign vs realne przerzuty) + raport Telegram (read-only)
Documentation=file:/root/.openclaw/workspace/scripts/dispatch_v2/tools/reassignment_shadow_eval.py
After=network.target
OnFailure=dispatch-onfailure-alert@%n.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/.openclaw/workspace/scripts
ExecStart=/root/.openclaw/venvs/dispatch/bin/python -m dispatch_v2.tools.reassignment_shadow_eval
StandardOutput=append:/root/.openclaw/workspace/scripts/logs/reassignment_shadow_eval.log
StandardError=append:/root/.openclaw/workspace/scripts/logs/reassignment_shadow_eval.log

[Install]
WantedBy=multi-user.target
