# deploy_staging/etc/systemd/system/dispatch-cod-panel-ingest.service.d/cron_health_success.conf
# FALA1 watchdog-close (2026-07-02) — audyt 2.0: false-"failed" oneshot.
# Runs & exits 0 weekly (verified: Jun 29 06:32 "Deactivated successfully") but had no
# success recorder → ledger frozen at failed/last_success=null since 2026-06-15.
# ExecStopPost records success (helper gates on SERVICE_RESULT=success) so the ledger
# reflects reality and cron_health.default_threshold_for backfills its 192h stale net.
# 1:1 with dispatch-cod-weekly-preflight.service.d/cron_health_success.conf.
[Service]
ExecStopPost=/root/.openclaw/workspace/scripts/dispatch_v2/observability/record_oneshot_success.sh %n
