# reconciliation/systemd/dispatch-state-reconcile.service
[Unit]
Description=Ziomek State Reconciliation Worker (TASK 2 Część B 2026-05-04)
After=network.target dispatch-shadow.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/.openclaw/workspace/scripts
Environment="PYTHONPATH=/root/.openclaw/workspace/scripts"
EnvironmentFile=/root/.openclaw/workspace/.secrets/telegram.env
ExecStart=/root/.openclaw/venvs/dispatch/bin/python -m dispatch_v2.reconciliation.reconcile_worker
StandardOutput=append:/root/.openclaw/workspace/scripts/logs/reconcile.log
StandardError=append:/root/.openclaw/workspace/scripts/logs/reconcile.log
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
