# systemd/dispatch-retro-learning.service.d/telegram-onfailure.conf
# Drop-in: Telegram alert gdy dzienna pętla retro-learning padnie.
# Mirror konwencji z dispatch-tg-heartbeat.service / dispatch-monitor-419.service.d/onfailure.conf.
# Template dispatch-onfailure-alert@.service jest zainstalowany w /etc/systemd/system/.
#
# Apply (część instalacji unitu, NIE osobno):
#   sudo mkdir -p /etc/systemd/system/dispatch-retro-learning.service.d
#   sudo cp dispatch-retro-learning.service.d/telegram-onfailure.conf \
#           /etc/systemd/system/dispatch-retro-learning.service.d/telegram-onfailure.conf
#   sudo systemctl daemon-reload
#
# Skutek: gdy łańcuch 4 kroków zerwie się (dowolny ExecStart exit != 0),
# unit kończy fail → odpala alert na Telegram. To CELOWE — chcemy wiedzieć,
# że trend przestał się akumulować.

[Unit]
OnFailure=dispatch-onfailure-alert@%n.service
