# docs/deploy/prep-bias-shadow/dispatch-prep-bias-shadow-monitor.timer
[Unit]
Description=Daily C2 prep-bias shadow monitor — 05:00 UTC (po dispatch-faza7-kpi backfill 04:00)

[Timer]
OnCalendar=*-*-* 05:00:00 UTC
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
