# deploy_staging/etc/systemd/system/dispatch-flag-fingerprint-guard.timer
[Unit]
Description=Ziomek STRAZNIK flag-fingerprint — co 30 min (rekoncyliacja stanu flag 4 zrodel; podsumowanie tick pisane ZAWSZE, tez przy OK). Wychwytuje rozjazd zrodel / intermittent-cold flag-load zanim zapadnie na zlych flagach.

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=dispatch-flag-fingerprint-guard.service

[Install]
WantedBy=timers.target
