# systemd/dispatch-plan-recheck.service.d/committed-propagation.conf
# Adrian 2026-06-22 — tie-breaker punktualności odbioru committed bez regresji SLA dostaw.
# Re-sekwencer (plan_recheck) doklejał OrderSim.czas_kuriera_warsaw → aktywuje miękką
# karę N5 + okno frozen w route_simulator_v2; przyjmuje wariant świadomy committed TYLKO
# gdy sla_violations dostaw <= baseline (R6 35min chroniony). Replay 22.06: punktualność
# ≤+5 27.6%→32.1% (świeże 58.8%→64.7%), Δ SLA dostaw -1/+0 (zero regresji), 0 INFEASIBLE.
# Monitor: grep COMMITTED_TIEBREAK_ADOPT/REJECT w logs/plan_recheck.log.
# Rollback: rm ten plik + systemctl daemon-reload (next tick = baseline).
[Service]
Environment=ENABLE_PLAN_RECHECK_COMMITTED_PROPAGATION=1
