# systemd/dispatch-new-courier-watch.service.d/onfailure.conf
# OnFailure → Telegram alert (parytet z pozostałymi dispatch-* oneshot).
# Apply przy deployu:
#   sudo mkdir -p /etc/systemd/system/dispatch-new-courier-watch.service.d
#   sudo cp onfailure.conf /etc/systemd/system/dispatch-new-courier-watch.service.d/
#   sudo systemctl daemon-reload
[Unit]
OnFailure=dispatch-onfailure-alert@%n.service
