# systemd/dispatch-shadow.service.d/override.conf
[Service]
# 2026-05-28 #33 latency outlier fix (commit V3.27.7 TECH_DEBT #20):
# bg_refresh thread eliminates ~22min CSRF expiry blocking call z hot path.
# Pre-flip: 124 outliers >3000ms w shadow_decisions.jsonl (login refresh w hot path).
# Worst case: oid=476410 11603ms (9s = panel_client login OK call sync).
# Expected post-flip: outliers >5000ms drop dramatically (Pre-V3.27.7 p95 ~6465ms).
Environment=ENABLE_PANEL_BG_REFRESH=1
Environment=ENABLE_BUNDLE_DELIV_SPREAD_CAP=1
Environment=ENABLE_LGBM_SHADOW=1
Environment=ENABLE_LGBM_METRICS_READ=1
# Sprint OBJ F2 (2026-05-18): koszt SPAN (idle) w objective OR-Tools.
# Zastępuje ENABLE_V328_P3D1_IDLE_COST (P3-D1 retired — był wadliwy).
# coeff=1.0 skalibrowany sweepem obj_harness (commit obj-f2-span-cost-2026-05-18).
Environment=ENABLE_OBJ_SPAN_COST=1
Environment=OBJ_SPAN_COST_COEFF=1.0
# Rolling late-binding Faza 0 (2026-05-18): shadow_dispatcher zasila pending_pool
# (czysta obserwacja). pending_pool_sweeper ma flagę w swoim unicie osobno.
Environment=ENABLE_PENDING_POOL=1
# Sprint OBJ F3 / BUG-4 (2026-05-18): best_effort z kandydatem łamiącym hard R6
# o > próg → KOORD zamiast auto-PROPOSE (diagnoza 474297). Próg 20 min (default).
Environment=ENABLE_OBJ_F3_BEST_EFFORT_R6_KOORD=1
# Sprint OBJ F4 Krok 1 (2026-05-18): proxy pozycji kuriera no-gps — picked_up
# → pickup_coords (punkt realny) zamiast delivery_coords (gdzie dopiero dojedzie).
# Eliminuje bias macierzy odległości (diagnoza 474266). commit 7098fee.
Environment=ENABLE_F4_COURIER_POS_PICKUP_PROXY=1
# Sprint OBJ F4 Krok 2 (2026-05-19): interpolacja pozycji kuriera no-gps na nodze
# pickup→delivery (f = clamp(elapsed/eta_leg, 0, 1), OSRM przez cache). Ma
# pierwszeństwo nad K1 gdy obie ON; fail-soft → K1 → legacy. commit 4c2cae8.
# Werdykt K1 (#54) PASS 2026-05-19 21:00 UTC. Ten restart wciąga też td20 K2
# fix (commit 4de3238) — koniec haversine spam z anchor (0,0).
Environment=ENABLE_F4_COURIER_POS_INTERP=1
