# systemd/dispatch-monitor-419.service.d/onfailure.conf
# Drop-in dla każdego dispatch-*.service:
# Plik: /etc/systemd/system/<unit>.service.d/onfailure.conf
#
# Apply:
#   sudo mkdir -p /etc/systemd/system/<unit>.service.d
#   sudo cp onfailure_dropin.conf /etc/systemd/system/<unit>.service.d/onfailure.conf
#   sudo systemctl daemon-reload
#
# Bulk apply do 14 unitów: zob. install_drop_ins.sh

[Unit]
OnFailure=dispatch-onfailure-alert@%n.service
