# systemd/dispatch-monitor-419.service.d/resource_limits.conf
# Profile: oneshot batch (daily/weekly/Mon batch jobs)
# Services: dispatch-daily-accounting, dispatch-cod-weekly-preflight,
#           dispatch-r04-evaluator, dispatch-overrides-reset,
#           dispatch-monitor-419
#
# Higher mem ceiling bo gspread/ML batch może być cięższe. Long timeout
# dla daily/weekly jobs.
#
# Apply: drop-in jako /etc/systemd/system/<unit>.service.d/resource_limits.conf

[Service]
MemoryHigh=400M
MemoryMax=600M
CPUQuota=100%
TasksMax=256
TimeoutStartSec=600
OOMScoreAdjust=100
